# rootfs/etc/services.d/coredns/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start CoreDNS service
# ==============================================================================
exec coredns -conf /config/corefile